# vm/systemctl/crusoe-log-collector-native.service
[Unit]
Description=Crusoe NVIDIA Log Collector for VMs (Native)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10
WorkingDirectory=/opt/crusoe-log-collector
EnvironmentFile=/etc/crusoe/crusoe_watch_agent/.env
EnvironmentFile=/etc/crusoe/secrets/.monitoring-token

# Environment variables for log collector
Environment="LOG_OUTPUT_DIR=/var/log/nvidia-bug-reports"
Environment="API_ENABLED=true"

ExecStart=/usr/bin/python3 /opt/crusoe-log-collector/vm_log_collector_app.py

# Security settings
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
