# vm/systemctl/crusoe-nvidia-log-collector.service
[Unit]
Description=Crusoe NVIDIA Log Collector for VMs
Documentation=https://github.com/crusoecloud/crusoe-watch-agent
After=docker.service crusoe-dcgm-exporter.service
Requires=docker.service
Wants=crusoe-dcgm-exporter.service

[Service]
Type=simple
TimeoutStartSec=0
Restart=always
RestartSec=10
WorkingDirectory=/etc/crusoe/crusoe_watch_agent

# Pull latest image before starting
ExecStartPre=/usr/bin/docker compose -f /etc/crusoe/crusoe_watch_agent/docker-compose-nvidia-log-collector.yaml pull

# Start the container
ExecStart=/usr/bin/docker compose -f /etc/crusoe/crusoe_watch_agent/docker-compose-nvidia-log-collector.yaml up

# Stop the container
ExecStop=/usr/bin/docker compose -f /etc/crusoe/crusoe_watch_agent/docker-compose-nvidia-log-collector.yaml down

[Install]
WantedBy=multi-user.target
